# .config/systemd/user/yubikey-cache-state.service
[Unit]
Description=Cache Yubikey-state

[Service]
Type=oneshot
ExecStart=%h/.local/bin/yubikey-cache-state.sh
