# .config/systemd/user/yubikey-cache-state.timer
[Unit]
Description=Update cached Yubikey-state

[Timer]
OnBootSec=5
OnUnitActiveSec=60

[Install]
WantedBy=timers.target
