# root/etc/update-motd.d/01-screenfetch
#!/usr/bin/env bash

printf "\n"
screenfetch -d '-host'
printf "\n"
